# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[inet_http_server]
port=*:9001
username=admin
password=ewFUyLEV*a4zxP

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:comfy]
command=/bin/bash -lc "cd ${CODE_DIR:-/comfy} && pixi run python main.py ${COMFY_ARGS}"
directory=/
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:proxy]
command=/root/.bun/bin/bun run start
directory=${CODE_DIR:-/comfy}/cj/proxy
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="8190",COMFY_URL="http://127.0.0.1:8188",JWT_ISSUERS="${PROXY_JWT_ISSUERS}",JWT_AUDIENCES="${PROXY_JWT_ISSUERS}",JWT_ALG_ALLOWLIST="${PROXY_JWT_ALG_ALLOWLIST:-EdDSA}",JWT_JWKS_URL="${PROXY_JWT_JWKS_URL}"
